# runit/nginx/run
#! /bin/ash

exec nginx \
    -c /srv/nginx.conf \
    -g 'daemon off;'
